# services/asr-core.service
[Unit]
Description=ASRCore - Shared ASR model service
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/projects/asr_core/.venv/bin/python3 -m asr_core
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=30
Environment=PYTHONUNBUFFERED=1
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
Environment=HF_DATASETS_OFFLINE=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
